# supervisor.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'

[program:tomcat]
command=/usr/local/tomcat/bin/catalina.sh run
